# PHASE2/rapidsms.crontab
MAILTO="pivot@pivotaccess.com"

# every day at 8 am, check our reminders
5 8 * * * cd /home/rapidsms/rapidsms/ && /usr/bin/python -W ignore::DeprecationWarning rapidsms checkreminders